# roles/udev/files/90-backlight.rules
# /etc/udev/rules.d/
# Set `video` as the owning group for the `/sys/class/backlight/intel_backlight/brightness` file
RUN+="/bin/chgrp video /sys/class/backlight/intel_backlight/brightness"

# Give write permisssions to the owning group of the `brightness` file
RUN+="/bin/chmod g+w /sys/class/backlight/intel_backlight/brightness"


